# netbird/rootfs/etc/services.d/netbird/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.warning "NetBird service stopped; s6 will restart it if configured."
