# transmission/rootfs/etc/services.d/transmission/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Transmission
# ==============================================================================
bashio::log.info "Starting Transmission..."

if bashio::config.true 'openvpn_enabled'; then
    exec /usr/sbin/openvpn --config /etc/openvpn/config.ovpn --script-security 2 --up /etc/openvpn/up.sh --down /etc/openvpn/down.sh
else
    exec /usr/bin/transmission-daemon --foreground --config-dir /data/transmission
fi
