# spotweb/rootfs/etc/periodic/hourly/spotweb
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Runs the hourly task to update spots
# ==============================================================================

bashio::log.info "Running hourly cron job" &>/proc/1/fd/1

s6-setuidgid nginx php /app/retrieve.php &>/proc/1/fd/1
