# manyfold/rootfs/etc/s6-overlay/s6-rc.d/manyfold/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
exec /run.sh
