# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_log/run
#!/usr/bin/with-contenv bash
s6-setuidgid ${USER}
echo "Starting service: birdnet log"
/usr/local/bin/gotty --address localhost -p 8080 --path log --title-format "BirdNET-Pi Log" birdnet_log.sh
