# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_log/run
#!/bin/execlineb -P
with-contenv
s6-envdir /run/user/1000
s6-setuidgid ${USER}
/usr/local/bin/gotty --address localhost -p 8080 --path log --title-format "BirdNET-Pi Log" birdnet_log.sh

#!/bin/sh
exec /usr/local/bin/gotty --address localhost -p 8080 --path log --title-format "BirdNET-Pi Log" birdnet_log.sh
