# seerr/rootfs/etc/services.d/nginx/finish
#!/usr/bin/execlineb -S0
# ==============================================================================
# Take down the S6 supervision tree when Nginx fails
# ==============================================================================
if { test ${1} -ne 0 }
if { test ${1} -ne 256 }

s6-svscanctl -t /var/run/s6/services
