# seerr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================

# Wait for Seerr to become available
bashio::net.wait_for 5055 localhost 900

bashio::log.info "Starting NGinx..."

exec nginx
