# seerr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================

# Set variables
slug=seerr
port=5055

# Wait for seerr to become available
bashio::net.wait_for "$port" localhost 900

bashio::log.info "Starting NGINX..."
exec nginx
