# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_analysis/run
#!/bin/execlineb -P
with-contenv
s6-setuidgid ${USER}
echo "Starting service: birdnet analysis"
$PYTHON_VIRTUAL_ENV /usr/local/bin/birdnet_analysis.py
