# ente/rootfs/etc/services.d/03-web/run
#!/usr/bin/with-contenv bash
set -e
[[ "${DISABLE_WEB_UI,,}" == "true" ]] && exit 0
exec /usr/bin/ente-web
