# zzz_draft_birdnet-pi/archive/s6-overlay/s6-rc.d/svc-chart_viewer/run
#!/usr/bin/with-contenv bashio
echo "Starting service: chart viewer"
exec \
    s6-setuidgid "$USER" "$PYTHON_VIRTUAL_ENV /usr/local/bin/daily_plot.py --daemon --sleep 2"
