# birdnet-pi/rootfs/templates/throttlescript.cron
#birdnet
*/15 * * * * $USER /usr/local/bin/throttlescript.sh >/dev/null 2>&1
