# tor/rootfs/etc/s6-overlay/s6-rc.d/tor/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tor
# Runs the Tor daemon
# ==============================================================================
bashio::log.info "Starting Tor..."

# Run the Tor daemon
exec tor
