# calibre_web/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================

# Wait for transmission to become available
bashio::net.wait_for 8083 localhost 900

bashio::log.info "Starting NGinx..."

exec nginx
