# readarr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

# Wait for transmission to become available
bashio::net.wait_for 8787 localhost 900

# Set base url

CONFIG_LOCATION=/config/addons_config/readarr/config.xml

if bashio::config.has_value 'CONFIG_LOCATION'; then
  CONFIG_LOCATION="$(bashio::config 'CONFIG_LOCATION')"
  # Modify if it is a base directory
  if [ -f "$CONFIG_LOCATION" ]; then CONFIG_LOCATION="$(dirname $CONFIG_LOCATION)/config.xml"; fi
fi

while ! [ -f "$CONFIG_LOCATION" ];
do
    echo "Waiting for $CONFIG_LOCATION to be created"
    sleep 1
done

if ! grep -q "<UrlBase>/readarr</UrlBase>" "$CONFIG_LOCATION"; then
  bashio::log.warning "BaseUrl not set properly, restarting"
  sed -i "/UrlBase/d" "$CONFIG_LOCATION"
  sed -i "/<Config>/a <UrlBase>\/readarr<\/UrlBase>" "$CONFIG_LOCATION"
  bashio::addon.restart
fi

bashio::log.info "Starting NGinx..."
exec nginx
