# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-getty/run
#!/bin/execlineb -P

with-contenv
s6-setuidgid ${USER}
echo "Starting service: getty"
/sbin/agetty --autologin ${USER} --noclear %I $TERM