# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-chart_viewer/run
#!/usr/bin/with-contenv bash
s6-setuidgid ${USER}
echo "Starting service: chart viewer"
$PYTHON_VIRTUAL_ENV /usr/local/bin/daily_plot.py --daemon --sleep 2
