# ente/rootfs/etc/services.d/02-minio/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

# Persistent path
MINIO_DATA="/config/minio-data"
mkdir -p "${MINIO_DATA}"

exec /usr/local/bin/minio server "${MINIO_DATA}" --address ":3200" & \
bashio::log.info "Starting MinIO (3200)..."
