# ente/rootfs/etc/services.d/02-minio/run
#!/usr/bin/with-contenv bash
set -e
exec /usr/local/bin/minio server /data --address ":3200"
