# aurral/rootfs/etc/services.d/aurral/run
#!/usr/bin/with-contenv bashio
# ==============================================================
# Start Aurral with configuration from HA add-on options.
# ==============================================================

DOWNLOAD_FOLDER=$(bashio::config 'download_folder')

# Ensure host-side download directory exists
mkdir -p "${DOWNLOAD_FOLDER}"

# Aurral writes its database/config to /app/backend/data inside the container.
# We bind that to HA's addon_config directory (/config) which is always writable.
# The upstream entrypoint chowns /app/backend/data — this must be a real directory,
# not a symlink to a protected HA path, so we do NOT symlink it.
mkdir -p /config/data

# If /app/backend/data exists as a real dir (from image), replace with our persistent dir
if [ ! -L /app/backend/data ]; then
  rm -rf /app/backend/data
  ln -sf /config/data /app/backend/data
fi

bashio::log.info "Starting Aurral"
bashio::log.info "  Downloads : ${DOWNLOAD_FOLDER}"
bashio::log.info "  Data      : /config/data"

export DOWNLOAD_FOLDER="${DOWNLOAD_FOLDER}"

exec /usr/local/bin/docker-entrypoint.sh node /app/server.js
