# aurral/rootfs/etc/services.d/aurral/run
#!/usr/bin/with-contenv bashio
# ==============================================================
# Start Aurral with configuration from HA add-on options.
# ==============================================================

DOWNLOAD_FOLDER=$(bashio::config 'download_folder')

# Ensure host-side download directory exists
mkdir -p "${DOWNLOAD_FOLDER}"

bashio::log.info "Starting Aurral"
bashio::log.info "  Data      : /data (HA persistent storage)"
bashio::log.info "  Downloads : ${DOWNLOAD_FOLDER}"

# /data is HA's private persistent storage, automatically mounted per-addon.
# AURRAL_DATA_DIR redirects Aurral's database/config there.
# DOWNLOAD_FOLDER is user-configured and points to /share.
export AURRAL_DATA_DIR="/data"
export DOWNLOAD_FOLDER="${DOWNLOAD_FOLDER}"
export WEEKLY_FLOW_FOLDER="${DOWNLOAD_FOLDER}/weekly-flow"

exec node /app/backend/server.js
