# aurral/rootfs/etc/services.d/aurral/run
#!/usr/bin/with-contenv bashio
# ==============================================================
# Start Aurral with configuration from HA add-on options.
# ==============================================================

DOWNLOAD_FOLDER=$(bashio::config 'download_folder')

# Ensure persistent directories exist
mkdir -p "${DOWNLOAD_FOLDER}"
mkdir -p /config/data
mkdir -p /config/data/image-proxy

# Pre-create /app/backend/data as a real directory so the upstream
# docker-entrypoint.sh can chown it without error.
mkdir -p /app/backend/data
ln -sfn /config/data/image-proxy /app/backend/data/image-proxy

bashio::log.info "Starting Aurral"
bashio::log.info "  Data      : /config/data (HA addon_config)"
bashio::log.info "  Downloads : ${DOWNLOAD_FOLDER}"

export PUID=0
export PGID=0
export AURRAL_DATA_DIR="/config/data"
export DOWNLOAD_FOLDER="${DOWNLOAD_FOLDER}"
export WEEKLY_FLOW_FOLDER="${DOWNLOAD_FOLDER}/weekly-flow"

exec docker-entrypoint.sh node server.js
