# aurral/rootfs/etc/services.d/aurral/run
#!/usr/bin/with-contenv bashio
# ==============================================================
# Start Aurral with configuration from HA add-on options.
# ==============================================================

DOWNLOAD_FOLDER=$(bashio::config 'download_folder')
DATA_FOLDER=$(bashio::config 'data_folder')

# Ensure host-side directories exist
mkdir -p "${DOWNLOAD_FOLDER}" "${DATA_FOLDER}"

bashio::log.info "Starting Aurral"
bashio::log.info "  Data      : ${DATA_FOLDER}"
bashio::log.info "  Downloads : ${DOWNLOAD_FOLDER}"

# Use Aurral's native env vars to redirect paths to HA share locations.
# This avoids symlinks into host-mounted paths which the upstream
# docker-entrypoint.sh cannot chown (Operation not permitted).
export AURRAL_DATA_DIR="${DATA_FOLDER}"
export DOWNLOAD_FOLDER="${DOWNLOAD_FOLDER}"
export WEEKLY_FLOW_FOLDER="${DOWNLOAD_FOLDER}/weekly-flow"

exec node /app/backend/server.js
