# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-caddy/run
#!/bin/execlineb -P
with-contenv
s6-setuidgid ${USER}
echo "Starting service: caddy"
/usr/local/bin/caddy run --config /etc/caddy/Caddyfile