# picard/rootfs/etc/services.d/picard/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Picard
# ==============================================================================
bashio::log.info "Starting Picard..."

exec /usr/bin/picard --foreground --config-dir /data/picard
