# claude_desktop/rootfs/etc/s6-overlay/s6-rc.d/svc-headroom/run
#!/usr/bin/with-contenv bashio
# Headroom optimization proxy — local backend for Claude Desktop MCP and Claude Code.
declare port=8787
declare host=127.0.0.1

# The dashboard is unauthenticated. Keep it container-local by default and bind all
# interfaces only when the user explicitly opts in and maps port 8787.
if bashio::config.true 'expose_headroom_dashboard'; then
    host=0.0.0.0
fi

if bashio::config.true 'install_headroom' && command -v headroom > /dev/null 2>&1; then
    # Kompress (the ONNX compression engine) needs its model in the local HF cache: the
    # proxy's startup preload is deliberately cache-only, and the default HF cache lands
    # under ~/.cache, which the add-on points at tmpfs (/tmp/cache) — wiped on every
    # restart. Without a warm persistent cache the proxy runs forever in "deferred" mode
    # and records zero compression savings. Cache on persistent storage and pre-warm it
    # once (bounded so an offline install still starts the proxy, in pass-through mode,
    # and retries the download on the next start).
    export HF_HOME="${HOME}/.headroom/hf"
    mkdir -p "$HF_HOME"
    chown abc:abc "$HF_HOME" 2> /dev/null || true
    # headroom lives in the lsiopy virtualenv; plain `python3` does not see its packages,
    # so derive the venv interpreter from the resolved headroom entry script.
    headroom_python="$(dirname "$(readlink -f "$(command -v headroom)")")/python3"
    if ! s6-setuidgid abc "$headroom_python" -c \
        'from headroom.transforms.kompress_compressor import _load_kompress; _load_kompress(allow_download=False)' \
        > /dev/null 2>&1; then
        bashio::log.info "svc-headroom: pre-warming the Kompress model cache (one-time download)"
        timeout 300 s6-setuidgid abc "$headroom_python" -c \
            'from headroom.transforms.kompress_compressor import _load_kompress; _load_kompress()' \
            || bashio::log.warning "svc-headroom: Kompress pre-warm failed; compression stays deferred until the model can be downloaded"
    fi
    bashio::log.info "svc-headroom: starting local Headroom proxy on ${host}:${port}"
    exec s6-setuidgid abc headroom proxy --host "${host}" --port "${port}" --code-aware
fi

bashio::log.info "svc-headroom: install_headroom disabled or headroom not found; idling"
exec sleep infinity
