# sonarr/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

# Set variables
slug=sonarr
port=8989
CONFIG_LOCATION=/config/addons_config/"$slug"/config.xml

# Wait for transmission to become available
bashio::net.wait_for "$port" localhost 900

if bashio::config.has_value 'CONFIG_LOCATION'; then
  CONFIG_LOCATION="$(bashio::config 'CONFIG_LOCATION')"
  # Modify if it is a base directory
  if [[ "$CONFIG_LOCATION" == *.* ]]; then CONFIG_LOCATION="$(dirname $CONFIG_LOCATION)"; fi
  CONFIG_LOCATION="$CONFIG_LOCATION"/config.xml
fi

# Set UrlBase
if ! grep -q "<UrlBase>/$slug</UrlBase>" "$CONFIG_LOCATION" && ! bashio::config.true "ingress_disabled"; then
  bashio::log.warning "BaseUrl not set properly, restarting"
  sed -i "/UrlBase/d" "$CONFIG_LOCATION"
  sed -i "/<Config>/a <UrlBase>\/$slug<\/UrlBase>" "$CONFIG_LOCATION"
  bashio::addon.restart
fi

# Disable Ingress
if bashio::config.true "ingress_disabled"; then
    bashio::log.info "Ingress is disabled, removing base url"

    # Remove BaseUrl
    if grep -q "$slug" "$CONFIG_LOCATION"; then
      bashio::log.warning "BaseUrl removed, restarting"
      sed -i "/$slug/d" "$CONFIG_LOCATION"
      bashio::addon.restart
    fi
fi

bashio::log.info "Starting NGinx..."
exec nginx
