# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-livestream/run
#!/bin/execlineb -P
with-contenv
echo "Starting service: livestream"
/usr/local/bin/livestream.sh