# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-livestream/run
#!/usr/bin/with-contenv bash
echo "Starting service: livestream"
/usr/local/bin/livestream.sh