# qbittorrent/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================

# Wait for transmission to become available
bashio::net.wait_for 8080 localhost 900

bashio::log.info "Starting NGinx..."

# Check vpn is working
if [ -f /currentip ]; then
    nginx || nginx -s reload &
    while true; do
        # Get vpn ip
        vpn_interface="tun0"
        use_interface=true

        if bashio::config.true 'wireguard_enabled'; then
            vpn_interface="${WIREGUARD_INTERFACE:-}"
            if [ -z "${vpn_interface}" ] && [ -f /var/run/s6/container_environment/WIREGUARD_INTERFACE ]; then
                vpn_interface="$(cat /var/run/s6/container_environment/WIREGUARD_INTERFACE)"
            fi
            vpn_interface="${vpn_interface:-wg0}"

            if bashio::config.true 'openvpn_alt_mode'; then
                use_interface=false
            fi
        else
            if bashio::config.true 'openvpn_alt_mode'; then
                use_interface=false
            fi
        fi

        if ${use_interface}; then
            if ! ip link show "${vpn_interface}" >/dev/null 2>&1; then
                bashio::log.warning "VPN interface ${vpn_interface} is not ready yet. Retrying in 5 seconds."
                sleep 5
                continue
            fi
            if ! curl -s ipecho.net/plain --interface "${vpn_interface}" > /vpnip; then
                bashio::log.warning "Unable to query public IP via interface ${vpn_interface}. Retrying in 5 seconds."
                sleep 5
                continue
            fi
        else
            if ! curl -s ipecho.net/plain > /vpnip; then
                bashio::log.warning "Unable to query public IP. Retrying in 5 seconds."
                sleep 5
                continue
            fi
        fi

        # Verify ip has changed
        if [[ "$(cat /vpnip)" = "$(cat /currentip)" ]]; then
            bashio::log.fatal "VPN is not properly configured. Your ip is exposed. Please fix this, or do not use the vpn alt mode"
            bashio::exit.nok
        fi

        # Get ip location
        COUNTRY=$(curl -s https://ipinfo.io/$(cat /vpnip) | grep country -i -m 1 | cut -d ':' -f 2 | xargs | awk 'gsub(/,$/,x)' || true)

        # Inform by message
        bashio::log.info "VPN is up and running with ip $(cat /vpnip), based in country : $COUNTRY"

        # Check every 15m
        sleep 15m

        true
    done
else
    nginx || nginx -s reload
fi
