# birdnet-pipy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail
exec nginx
