# zzz_draft_birdnet-pi/archive/s6-overlay/s6-rc.d/svc-livestream/run
#!/usr/bin/with-contenv bashio
echo "Starting service: livestream"
exec \
    s6-setuidgid "$USER" "/usr/local/bin/livestream.sh"
