# transmission_ls/rootfs/etc/services.d/transmission/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Transmission
# ==============================================================================
bashio::log.info "Starting Transmission..."

exec /usr/bin/transmission-daemon --foreground --config-dir /share/transmission