# spotweb/rootfs/etc/services.d/phpfpm/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Starting PHP-FPM..."

#/usr/sbin/php-fpm* -F
exec php-fpm* --nodaemonize
