# hyperion/root/etc/services.d/hyperion/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Starting hyperiond...'

exec hyperiond -u /config/hyperion