# portainer_agent/rootfs/etc/services.d/portainer_agent/finish
#!/usr/bin/execlineb -S1

s6-svscanctl -t /var/run/s6/services