# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_recording/finish
#!/usr/bin/with-contenv bash

if { sleep 6 }
s6-svc -u /etc/s6-overlay/s6-rc.d/svc-birdnet_recording