# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_recording/run
#!/usr/bin/with-contenv bash
s6-setuidgid ${USER}
echo "Starting service: birdnet recording"
/usr/local/bin/birdnet_recording.sh