# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-birdnet_recording/run
#!/bin/execlineb -P
with-contenv
s6-setuidgid ${USER}
/usr/local/bin/birdnet_recording.sh