# netbird-server/rootfs/etc/services.d/signal/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

# ==============================================================================
# Home Assistant Add-on: NetBird Server
# Runs the NetBird Signal service
# ==============================================================================

SIGNAL_PORT=8083
LOG_LEVEL="info"

bashio::log.info "Starting NetBird Signal on port ${SIGNAL_PORT} (legacy gRPC on 10000)..."
exec /usr/local/bin/netbird-signal run \
  --port "$SIGNAL_PORT" \
  --log-level "$LOG_LEVEL" \
  --log-file console
