# ente/rootfs/etc/services.d/01-minio/run
#!/usr/bin/with-contenv bashio
set -euo pipefail
bashio::log.info "Starting MinIO."
exec /usr/local/bin/minio server /data --address ":3200"
