# bitwarden/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Vaultwarden
# Runs the Nginx daemon
# ==============================================================================
bashio::net.wait_for 80
bashio::log.info "Starting NGinx..."

exec nginx
