# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/svc-qbittorrent/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

WEBUI_PORT=${WEBUI_PORT:-8080}
export PATH="/usr/local/sbin:/usr/local/bin:${PATH}"

# --- Configuration & Pre-checks ---

if bashio::config.true 'silent'; then
    sed -i 's|/proc/1/fd/1 hassio;|off;|g' /etc/nginx/nginx.conf
fi

# --- Main Execution ---

openvpn_enabled=false
wireguard_enabled=false

if bashio::config.true 'openvpn_enabled'; then
    openvpn_enabled=true
fi

if bashio::config.true 'wireguard_enabled'; then
    wireguard_enabled=true
fi

if [[ "${openvpn_enabled}" == true && "${wireguard_enabled}" == true ]]; then
    bashio::log.error "Both OpenVPN and WireGuard are enabled. Only one VPN mode is supported."
    exit 1
fi

if [[ "${openvpn_enabled}" == true ]]; then
    /usr/local/sbin/vpn openvpn up
elif [[ "${wireguard_enabled}" == true ]]; then
    /usr/local/sbin/vpn wireguard up
fi

# --- Launch qBittorrent ---

bashio::log.info "Starting qBittorrent..."

# On aarch64 LSIO images svc-qbittorrent has no notification-fd, so fd 3 is never
# opened by s6-rc. s6-notifyoncheck then exits with EBADF after its check passes,
# s6 sees the supervised process die and restarts it — causing the startup loop.
# Run qbittorrent-nox directly so s6 supervises the actual process.
if bashio::config.true 'silent'; then
    exec > /dev/null 2>&1
fi

exec s6-setuidgid abc /app/qbittorrent-nox --webui-port="${WEBUI_PORT}"
