# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/svc-qbittorrent/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

WEBUI_PORT=${WEBUI_PORT:-8080}
export PATH="/usr/local/sbin:/usr/local/bin:${PATH}"

# Read WebUI bind address from config (mirrors upstream LSIO approach)
WEBUI_ADDRESS=$(grep -Po "^WebUI\\\\Address=\K(.*)" /config/qBittorrent/qBittorrent.conf 2>/dev/null || echo "")
if [[ -z ${WEBUI_ADDRESS} ]] || [[ ${WEBUI_ADDRESS} == "*" ]]; then
    WEBUI_ADDRESS="localhost"
fi

# --- Configuration & Pre-checks ---

if bashio::config.true 'silent'; then
    sed -i 's|/proc/1/fd/1 hassio;|off;|g' /etc/nginx/nginx.conf
fi

# --- Main Execution ---

openvpn_enabled=false
wireguard_enabled=false

if bashio::config.true 'openvpn_enabled'; then
    openvpn_enabled=true
fi

if bashio::config.true 'wireguard_enabled'; then
    wireguard_enabled=true
fi

if [[ "${openvpn_enabled}" == true && "${wireguard_enabled}" == true ]]; then
    bashio::log.error "Both OpenVPN and WireGuard are enabled. Only one VPN mode is supported."
    exit 1
fi

if [[ "${openvpn_enabled}" == true ]]; then
    /usr/local/sbin/vpn openvpn up
elif [[ "${wireguard_enabled}" == true ]]; then
    /usr/local/sbin/vpn wireguard up
fi

# --- Launch qBittorrent ---

# Find qbittorrent-nox binary across different LSIO image layouts
QB_BIN=""
for candidate in /app/qbittorrent-nox /usr/bin/qbittorrent-nox /usr/local/bin/qbittorrent-nox; do
    if [[ -x "${candidate}" ]]; then
        QB_BIN="${candidate}"
        break
    fi
done

if [[ -z "${QB_BIN}" ]]; then
    bashio::log.fatal "qbittorrent-nox binary not found"
    exit 1
fi

# Find qbittorrent-nox binary across different LSIO image layouts
QB_BIN=""
for candidate in /app/qbittorrent-nox /usr/bin/qbittorrent-nox /usr/local/bin/qbittorrent-nox; do
    if [[ -x "${candidate}" ]]; then
        QB_BIN="${candidate}"
        break
    fi
done

if [[ -z "${QB_BIN}" ]]; then
    bashio::log.fatal "qbittorrent-nox binary not found"
    exit 1
fi

bashio::log.info "Starting qBittorrent..."

# Suppress output in silent mode by redirecting the current shell's stdout/stderr
# before exec replaces the process (avoids opening /dev/stdout which may fail on
# some ARM platforms / container runtimes).
if bashio::config.true 'silent'; then
    exec >/dev/null 2>&1
fi

if [[ -z ${LSIO_NON_ROOT_USER} ]]; then
    exec s6-setuidgid abc "${QB_BIN}" --webui-port="${WEBUI_PORT}"
else
    exec "${QB_BIN}" --webui-port="${WEBUI_PORT}"
fi
