# qbittorrent/rootfs/etc/s6-overlay/s6-rc.d/svc-qbittorrent/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

WEBUI_PORT=${WEBUI_PORT:-8080}

export PATH="/usr/local/sbin:/usr/local/bin:${PATH}"

wireguard_firewall_command_works() {
    local cmd="$1"

    if ! command -v "${cmd}" >/dev/null 2>&1; then
        return 1
    fi

    if ! { printf '%s\n' '*filter'; printf '%s\n' 'COMMIT'; } | "${cmd}" -n >/dev/null 2>&1; then
        return 1
    fi

    return 0
}

wireguard_create_firewall_wrapper() {
    local cmd="$1"
    local dir="$2"
    local target="${dir}/${cmd}"

    {
        echo '#!/bin/sh'
        printf 'echo "WireGuard firewall fallback: %s skipped." >&2\n' "${cmd}"
        echo 'exit 0'
    } > "${target}"

    chmod +x "${target}"
}

wireguard_prepare_firewall_wrappers() {
    local ipv6_state="$1"
    local fallback_dir="${WIREGUARD_STATE_DIR}/iptables-fallback"
    local wrappers_created=false

    rm -rf "${fallback_dir}"
    mkdir -p "${fallback_dir}"

    if ! wireguard_firewall_command_works 'iptables-restore'; then
        wrappers_created=true
        bashio::log.warning 'iptables-restore is unavailable or failed to initialize. WireGuard will skip iptables firewall programming.'
        wireguard_create_firewall_wrapper 'iptables-restore' "${fallback_dir}"
    fi

    if [[ "${ipv6_state}" != 'disabled' ]] && ! wireguard_firewall_command_works 'ip6tables-restore'; then
        wrappers_created=true
        bashio::log.warning 'ip6tables-restore is unavailable or failed to initialize. WireGuard will skip IPv6 firewall programming.'
        wireguard_create_firewall_wrapper 'ip6tables-restore' "${fallback_dir}"
    fi

    if [[ "${wrappers_created}" == true ]]; then
        export PATH="${fallback_dir}:${PATH}"
        bashio::log.warning 'Continuing without kernel iptables support. Consider enabling the required modules on the host for full kill-switch protection.'
    else
        rmdir "${fallback_dir}" 2>/dev/null || true
    fi
}

if bashio::config.true 'silent'; then
    sed -i 's|/proc/1/fd/1 hassio;|off;|g' /etc/nginx/nginx.conf
fi

if bashio::config.true 'openvpn_enabled'; then
    exec /usr/sbin/openvpn \
        --config /config/openvpn/config.ovpn \
        --script-security 2 \
        --up /etc/openvpn/up.sh \
        --down /etc/openvpn/down.sh \
        --pull-filter ignore "route-ipv6" \
        --pull-filter ignore "ifconfig-ipv6" \
        --pull-filter ignore "tun-ipv6" \
        --pull-filter ignore "redirect-gateway ipv6" \
        --pull-filter ignore "dhcp-option DNS6"
else
    if bashio::config.true 'wireguard_enabled'; then
        WIREGUARD_STATE_DIR="/var/run/wireguard"

        if ! bashio::fs.file_exists "${WIREGUARD_STATE_DIR}/config"; then
            bashio::exit.nok 'WireGuard runtime configuration not prepared. Please restart the add-on.'
        fi

        wireguard_config="$(cat "${WIREGUARD_STATE_DIR}/config")"
        wireguard_interface="$(cat "${WIREGUARD_STATE_DIR}/interface" 2>/dev/null || echo 'wg0')"
        wireguard_ipv6_state="$(cat "${WIREGUARD_STATE_DIR}/ipv6_state" 2>/dev/null || echo 'unknown')"

        wireguard_prepare_firewall_wrappers "${wireguard_ipv6_state}"

        if ip link show "${wireguard_interface}" &> /dev/null; then
            bashio::log.warning "WireGuard interface ${wireguard_interface} already exists. Attempting to reset it."
            wg-quick down "${wireguard_config}" >/dev/null 2>&1 || true
        fi

        bashio::log.info "Starting WireGuard interface ${wireguard_interface} using ${wireguard_config##*/}."

        if ! output=$(wg-quick up "${wireguard_config}" 2>&1); then
            bashio::log.error 'WireGuard failed to establish a connection.'
            bashio::log.error "wg-quick output:"
            bashio::log.error "${output}"
            bashio::log.error 'Troubleshooting steps:'
            bashio::log.error "  1. Confirm that the WireGuard configuration file '${wireguard_config}' exists inside the container and contains valid private/public keys, endpoint and AllowedIPs."
            bashio::log.error '  2. Ensure UDP port 51820 (or the port defined in your config) is forwarded on your router to this host and not blocked by your firewall or ISP.'
            bashio::log.error '  3. Verify that the configured endpoint (IP/hostname and port) is reachable from this container (e.g. ping or nc from a debug shell).'
            bashio::log.error '  4. Check that the system time is correct (NTP); large time drift can break key handshakes.'
            bashio::log.error '  5. Confirm that WireGuard kernel support / module is available in the host system.'
            bashio::log.error '  6. If DNS names are used for the endpoint, verify DNS resolution from inside the container (e.g. nslookup or dig).'
            bashio::exit.nok 'WireGuard start failed. See the log above for details.'
        fi

        bashio::log.info "WireGuard interface ${wireguard_interface} is up."

        # Refresh DNS resolver configuration if resolvconf is present
        if command -v resolvconf >/dev/null 2>&1; then
            bashio::log.info 'Refreshing DNS resolver configuration via resolvconf -u.'
            if ! resolvconf -u >/dev/null 2>&1; then
                bashio::log.warning 'resolvconf -u failed. DNS configuration may not have been updated.'
            fi
        else
            bashio::log.debug 'resolvconf not found in PATH; skipping DNS refresh.'
        fi
    fi

    if bashio::config.true 'silent'; then
        exec \
            s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${WEBUI_PORT}" \
            s6-setuidgid abc /usr/bin/qbittorrent-nox --webui-port="${WEBUI_PORT}" > /dev/null
    else
        exec \
            s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${WEBUI_PORT}" \
            s6-setuidgid abc /usr/bin/qbittorrent-nox --webui-port="${WEBUI_PORT}"
    fi
fi
