# birdnet-pipy/rootfs/etc/services.d/icecast/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail
export PULSE_SERVER="${PULSE_SERVER:-unix:/run/pulse/native}"
gosu icecast2 /usr/local/bin/start-icecast.sh
