# transmission-magnet-redirect/root/etc/services.d/transmission-magnet-redirect/run
#!/usr/bin/with-contenv bashio

bashio::log.info 'Starting transmission-magnet-redirect...'

cd /app
exec dotnet TransmissionMagnetRedirect.dll