# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-spectrogram_viewer/finish
#!/usr/bin/with-contenv bash

if { sleep 10 }
s6-svc -u /etc/s6-overlay/s6-rc.d/svc-spectrogram_viewer