# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-spectrogram_viewer/run
#!/bin/execlineb -P
with-contenv
s6-setuidgid ${USER}
/usr/local/bin/spectrogram.sh
