# karakeep/rootfs/etc/s6-overlay/s6-rc.d/svc-meilisearch/run
#!/command/with-contenv bashio
# shellcheck shell=bash
set -e

mkdir -p "${MEILI_DIR}"

exec /bin/meilisearch --db-path "${MEILI_DIR}" --no-analytics --experimental-dumpless-upgrade
