# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-avahi-alias/run
#!/usr/bin/with-contenv bashio
s6-setuidgid ${USER}

echo "Starting service: avahi alias"
/bin/bash -c "/usr/bin/avahi-publish -a -R %I $(hostname -I |cut -d' ' -f1)"