# netbird-server/rootfs/etc/services.d/coturn/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

# ==============================================================================
# Home Assistant Add-on: NetBird Server
# Runs Coturn
# ==============================================================================

DATA_DIR="/config/netbird"
TURN_CONFIG="$DATA_DIR/turn/turnserver.conf"

if [[ ! -f "$TURN_CONFIG" ]]; then
  bashio::log.error "Missing Coturn configuration at ${TURN_CONFIG}."
  bashio::exit.nok
fi

bashio::log.info "Starting Coturn..."
exec /usr/local/bin/turnserver -c "$TURN_CONFIG" --log-file stdout
