# birdnet-pi/rootfs/scripts/throttlescript.cron
#birdnet
*/15 * * * * $USER /scripts/throttlescript.sh >/dev/null 2>&1
