# nextcloud/rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

# log errors in addon logs in foreground
/usr/sbin/crond -f -L /proc/1/fd/1 -l 3 -c /etc/crontabs
