# portainer/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================

bashio::log.info "Waiting for port 9000 to open..."

# Wait for transmission to become available
bashio::net.wait_for 9000 localhost 900

bashio::log.info "Starting NGinx..."

exec nginx
