# ente/rootfs/etc/services.d/05-web/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

if bashio::config.true 'DISABLE_WEB_UI'; then
    bashio::log.info "Web UI disabled."
    exit 0
fi

exec /usr/bin/ente-web & \
bashio::log.info "Starting Ente web."
