# bitwarden/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e
# ==============================================================================
# Home Assistant Community Add-on: Bitwarden
# Runs the Nginx daemon
# ==============================================================================
bashio::net.wait_for 80
bashio::log.info "Starting NGinx..."

exec nginx
