# zzz_draft_birdnet-pi/rootfs/etc/s6-overlay/s6-rc.d/svc-web_terminal/run
#!/bin/execlineb -P
with-contenv
/usr/local/bin/gotty --address localhost -w -p 8888 --path terminal --title-format "BirdNET-Pi Terminal" login
