# vaultwarden/rootfs/etc/services.d/bitwarden/finish
#!/command/with-contenv bashio 
# ============================================================================== 
# Home Assistant Community Add-on: Vaultwarden 
# Take down the S6 supervision tree when Nginx fails 
# ============================================================================== 

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then 
bashio::log.warning "NGINX crashed, halting add-on" 
/run/s6/basedir/bin/halt 
fi 

bashio::log.info "NGINX stopped, restarting..." 
