# zzz_draft_birdnet-pi2/archive/s6-overlay/s6-rc.d/svc-birdnet_log/run
#!/usr/bin/with-contenv bashio
echo "Starting service: birdnet log"
exec \
    s6-setuidgid "$USER" "/usr/local/bin/gotty --address localhost -p 8080 --path log --title-format 'BirdNET-Pi Log' birdnet_log.sh"
