# zzz_draft_birdnet-pi/archive/s6-overlay/s6-rc.d/svc-birdnet_recording/run
#!/usr/bin/with-contenv bashio
echo "Starting service: birdnet recording"
exec \
    s6-setuidgid "$USER" "/usr/local/bin/birdnet_recording.sh"
