# ente/rootfs/etc/services.d/03-minio-init/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

MINIO_USER="$(bashio::config 'MINIO_ROOT_USER')"
MINIO_PASS="$(bashio::config 'MINIO_ROOT_PASSWORD')"
S3_BUCKET="$(bashio::config 'S3_BUCKET')"

bashio::log.info "MinIO-init: waiting for API..."
until /usr/local/bin/mc alias set h0 http://127.0.0.1:3200 "${MINIO_USER}" "${MINIO_PASS}" 2>/dev/null; do
    sleep 1
done

bashio::log.info "MinIO-init: ensuring bucket ${S3_BUCKET}..."
/usr/local/bin/mc mb -p "h0/${S3_BUCKET}" || true

bashio::log.info "MinIO-init: done."
exit 0
