# bentopdf/rootfs/etc/s6-overlay/s6-rc.d/bentopdf/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash
exit 0
